# ENTER/etc/profile.d/conda.csh
setenv CONDA_EXE "/mnt/sevenT/spinfer/SpInfer_EuroSys25/ENTER/bin/conda";
setenv _CONDA_ROOT "/mnt/sevenT/spinfer/SpInfer_EuroSys25/ENTER";
setenv _CONDA_EXE "/mnt/sevenT/spinfer/SpInfer_EuroSys25/ENTER/bin/conda";
setenv CONDA_PYTHON_EXE "/mnt/sevenT/spinfer/SpInfer_EuroSys25/ENTER/bin/python";
source "/mnt/sevenT/spinfer/SpInfer_EuroSys25/ENTER/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh";
